# cron jobs/cron.weekly/schedulePurge.sh
#!/bin/bash

#Author: Robert Gordon
#Date: 29 Jan 2018

echo "Scheduling purge.py..."
echo 'Running purge.py...'
/Users/robudhabi/purge.py -v
echo "Finshed purging the system."
