# cron jobs/cron.daily/scheduleDiskCheck.sh
#!/bin/bash

#Author: Robert Gordon
#Date: 1 Jan 2018

echo "Scheduling disk space check.."
echo "Executing 'checkDiskSpace.py'..."
/Users/robudhabi/checkDiskSpace.py
echo "Finshed checking disk."
